# root/etc/services.d/mkdocs/run
#!/usr/bin/with-contenv bash
#cd /config
#mkdocs serve
