# root/etc/services.d/mkdocs/run
#!/usr/bin/with-contenv bash
cd /config
exec \
	s6-setuidgid abc mkdocs serve --dev-addr 0.0.0.0:8000
